# loc_hack/src/debrick.conf
# Script run after all filesystems have loaded

start on started filesystems_userstore
stop on stopping filesystems

export LANG LC_ALL

pre-start script
  /etc/upstart/debrick &
end script
